# ScoutTests/Fixtures/plists/com.scout.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scout.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/scout-dev/Scout/scripts/heartbeat.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/scout-dev/Scout</string>

    <!-- Run every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Environment: ensure PATH includes homebrew, gh, and claude -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/scout-dev/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/scout-dev</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/scout-dev/Scout/.scout-logs/launchd-heartbeat-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/scout-dev/Scout/.scout-logs/launchd-heartbeat-stderr.log</string>

</dict>
</plist>
