# ScoutTests/Fixtures/plists/com.scout.dreaming.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scout.dreaming</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/scout-dev/Scout/run-dreaming.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/scout-dev/Scout</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 6:33 PM Mon-Fri (Dreaming) -->
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
        <!-- 8:33 PM Mon-Fri (Dreaming) -->
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>33</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <!-- Environment: ensure PATH includes homebrew, gh, and claude -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/scout-dev/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/scout-dev</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/scout-dev/Scout/.scout-logs/launchd-dreaming-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/scout-dev/Scout/.scout-logs/launchd-dreaming-stderr.log</string>

</dict>
</plist>
