# ScoutTests/Fixtures/plists/com.scout.briefing-weekend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scout.briefing-weekend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/scout-dev/Scout/run-scout.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <!-- Saturday 8:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Sunday 8:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/scout-dev/Scout/.scout-logs/launchd-briefing-weekend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/scout-dev/Scout/.scout-logs/launchd-briefing-weekend.log</string>
</dict>
</plist>
